# systemd-home/.config/systemd/user/notifi.service
[Unit]
Description=start notifi script

[Service]
ExecStart=/home/vlkn/grepo/bin/notifi
Restart=always

[Install]
WantedBy=multi-user.target
