# systemd-home/.config/systemd/user/pipgem.service
[Unit]
Description=update pip packages

[Service]
Type=oneshot
ExecStart=/home/vlkn/.local/bin/pipgem
NoNewPrivileges=true
