# systemd-etc/etc/systemd/system/ccnow.service
[Unit]
Description=cleans RAM and swap cache
AssertPathExists=/usr/local/bin/ccnow

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ccnow
