# systemd-etc/etc/systemd/system/mirrorlist.timer
[Unit]
Description=weekly generate mirrorlist

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
