# etc/systemd/btrfs.service
[Unit]
Description=btrfs scrub on %f

[Service]
Nice=19
IOSchedulingClass=Idle
ExecStart=/usr/bin/btrfs scrub start -B %f
