# etc/systemd/btrfs.timer
[Unit]
Description=Monthly btrfs scrub on %f

[Timer]
OnCalendar=monthly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=multi-user.target
