# etc/systemd/ccnow.timer
[Unit]
Description=runs cache cleaning every 3hrs

[Timer]
OnBootSec=3h
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
